# ubuntu/log2kafka.upstart
description  "forward logs from file to Kafka"
author       "Martin Schuette <martin.schuette@deck36.de>"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
console log # log stdout/stderr to /var/log/upstart/

script
  . /etc/default/log2kafka
  exec start-stop-daemon --start -c $SETUID --exec /usr/local/bin/log2kafka -- -file="$LOGFILE" -server="$BROKER" -topic="$TOPIC"
end script
